# config/etc/profile.d/ruby_enterprise_path.sh
PATH=$PATH:/opt/ruby-enterprise-1.8.7-2010.02/bin
export PATH
