# config/etc/profile.d/ruby_enterprise_path.sh
PATH=$PATH:/opt/ruby-enterprise-1.8.6-20090610/bin
export PATH
